# scripts/systemd/llm-quant-commodities.timer
[Unit]
Description=Run llm-quant commodities every 5 minutes

[Timer]
OnCalendar=Mon..Fri *-*-* *:0/5:00
Persistent=true

[Install]
WantedBy=timers.target
