# scripts/systemd/llm-quant-crypto.service
[Unit]
Description=llm-quant crypto intraday run

[Service]
Type=oneshot
WorkingDirectory=/home/ty/Documents/llm-quant/llm-quant
EnvironmentFile=/home/ty/Documents/llm-quant/llm-quant/.env
ExecStart=/home/ty/Documents/llm-quant/llm-quant/.venv/bin/pq run --pod crypto --broker alpaca
