# scripts/systemd/llm-quant-eod-flat.timer
[Unit]
Description=llm-quant end-of-day flatten timer

[Timer]
OnCalendar=Mon..Fri *-*-* 15:55:00
Timezone=America/New_York
Persistent=true

[Install]
WantedBy=timers.target
