# scripts/systemd/llm-quant-eod-flat-stock-bounded-paper.timer
[Unit]
Description=llm-quant end-of-day flatten timer (stock-bounded-paper)

[Timer]
# User-systemd on some hosts ignores Timezone=, so this schedule follows the
# user manager's local timezone. Run under an America/New_York user/host
# timezone when exact 15:55 ET semantics are required.
OnCalendar=Mon..Fri *-*-* 15:55:00
Persistent=true
Unit=llm-quant-eod-flat-stock-bounded-paper.service

[Install]
WantedBy=timers.target
